# ps-startup.service
[Unit]
After=network.service

[Service]
ExecStart=/home/pi/hub/execpipe.sh

[Install]
WantedBy=default.target
